# backend/systemd-storyapp.service
[Unit]
Description=Story App (Node backend + static frontend)
Wants=network-online.target sound.target
After=network-online.target sound.target

[Service]
Type=simple
User=storyapp
Group=storyapp
EnvironmentFile=/etc/storyapp/env
WorkingDirectory=/opt/storyapp/current
ExecStart=/bin/sh -lc 'exec node backend/server.js >> /var/log/storyapp/app.log 2>&1'
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ReadWritePaths=/var/log/storyapp /var/lib/storyapp /opt/storyapp

[Install]
WantedBy=multi-user.target
