# deploy/storyapp.service
# Minimal production systemd unit for Pi Zero 2W
[Unit]
Description=Bedtime Stories App
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/storyapp/backend
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
