# deploy/storyapp.service
## Minimal production systemd unit (simplified) ##
[Unit]
Description=Bedtime Stories App
## Minimal production systemd unit (final simplified version) ##
[Unit]
Description=Bedtime Stories App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/storyapp/backend
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
# with comprehensive security, resource management, and monitoring capabilities.
